# klipper/mjpg_streamer.ini

[program:mjpg-streamer]
user=klippy
command=/home/klippy/mjpg-streamer/mjpg_streamer -o "output_http.so -w ./www-mjpgstreamer -n" -i "input_uvc.so -d /dev/video0 -"
environment=LD_LIBRARY_PATH=/home/klippy/mjpg-streamer
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes = 0
autorestart=true
autostart=false
dependent_startup=true
