# klipper/moonraker.ini
[program:moonraker]
user=klippy
command=run_in_venv /home/klippy/moonraker-env python /home/klippy/moonraker/moonraker/moonraker.py -l /var/log/klipper/moonraker.log -c /home/klippy/.config/moonraker.conf
environment=USER=klippy,HOME=/home/klippy
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes = 0
autorestart=true
autostart=false
dependent_startup=true
dependent_startup_wait_for=klipper:running
